# profile.d/Darwin
alias top='top -o cpu'
alias hide='/usr/bin/SetFile -a V'
export COPY_EXTENDED_ATTRIBUTES_DISABLE=true
export COPYFILE_DISABLE=true
export CLASSPATH="${HOME}/Library/Java"
export LEIN_HOME="${HOME}/Library/Leiningen"
export GEM_HOME="${HOME}/Library/Ruby"
export HAXE_STD_PATH="${HOME}/Library/Haxe/std"
export HOMEBREW_CASK_OPTS='--appdir=/Users/Shared/Applications --caskroom=/usr/local/Caskroom'
export N_PREFIX="${HOME}/Library/Node"
export PERL_MB_OPT="--install_base ${HOME}/Library/Perl";
export PERL_MM_OPT="INSTALL_BASE=${HOME}/Library/Perl";
export PERL5LIB="${HOME}/Library/Perl/lib/perl5";
export PYTHONPATH="${HOME}/Library/Python"
export VAGRANT_DEFAULT_PROVIDER="vmware_fusion"
export VAGRANT_HOME="${HOME}/Library/Vagrant"
