# profile.d/Darwin
alias top='top -o cpu'
alias hide='/usr/bin/SetFile -a V'
export COPY_EXTENDED_ATTRIBUTES_DISABLE=true
export COPYFILE_DISABLE=true
export CLASSPATH="${HOME}/Library/Java"
export LEIN_HOME="${HOME}/Library/Leiningen"
export GEM_HOME="${HOME}/Library/Ruby"
export GOPATH="${HOME}/Library/Go"
export HAXE_STD_PATH="${HOME}/Library/Haxe/std"
export PYTHONPATH="${HOME}/Library/Python"
PATH=$PATH:/usr/local/texlive/2017basic/bin/x86_64-darwin
