# bash/.config/profile.d/Darwin
alias top='top -o cpu'
alias hide='/usr/bin/SetFile -a V'
export COPY_EXTENDED_ATTRIBUTES_DISABLE=true
export COPYFILE_DISABLE=true
export CLASSPATH="${HOME}/Library/Java"
export GOPATH="${HOME}/Library/Go"
export PYTHONPATH="${HOME}/Library/Python"
export RBENV_ROOT="${HOME}/Library/Ruby"
PATH=/usr/local/texlive/2017basic/bin/x86_64-darwin:$PATH
