# bash/.config/profile.d/Linux
alias cal='cal -N -M'
alias vi='vim'

export GOPATH="${HOME}/lib/go"
export RBENV_ROOT="${HOME}/lib/ruby"

PATH=${HOME}/bin:$PATH

if [ -n "$PS1" ]; then
  _available ssh-agent && eval "$(ssh-agent -s)"
  trap '/usr/bin/ssh-agent -k' 0
fi
