# bash/.config/profile.d/Linux
alias vi='vim'

export RBENV_ROOT="${HOME}/lib/ruby"

PATH=${HOME}/bin:$PATH
